# vsysmon.service
[Unit]
Description=Vsysmon Daemon
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/vsysmon-daemon -p 50051 -c /etc/vsysmon/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vsysmon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target